# debian/postinst
#!/bin/sh

set -e 

update-rc.d nginx defaults >/dev/null
service nginx start

exit 0
