# debian/postinst
#!/bin/sh

set -e 

update-rc.d wisdom-media-server defaults >/dev/null
service wisdom-media-server start

exit 0
