# debian/postrm
#!/bin/sh

set -e

update-rc.d wisdom-media-server remove >/dev/null
rm -rf /opt/wisdom-media-server

exit 0
