# debian/postrm
#!/bin/sh

set -e

update-rc.d nginx remove >/dev/null
rm -rf /opt/nginx

exit 0
